# ansible/roles/etc-environment/files/etc-environment.target
[Unit]
Description=Seed /etc/environment from cloud-config
Requires=etc-environment.service
After=etc-environment.service
